# admin/boot/gentoo/openexam-render
#!/sbin/openrc-run
# 
# Copyright (C) 2017 The OpenExam project.
# Distributed under the terms of the GNU General Public License, v2 or later

description="openexam render worker"
pidfile="/var/run/openexam/render.pid"

command="/usr/local/bin/openexam-render"
command_args="--start"
command_user="apache"

required_dirs="/var/run/openexam"
required_files="/usr/local/bin/openexam-render"

depend() {
        need apache2
}

start_pre() {
        mkdir -p /var/run/openexam
        cd /
}

start() {
        ebegin "Starting OpenExam Render Worker"
        start-stop-daemon --start \
            --make-pidfile --pidfile ${pidfile} \
            --exec ${command} -- ${command_args} >& /dev/null
        eend $?
}

stop() {
        ebegin "Starting OpenExam Render Worker"
        ${command} --stop >& /dev/null
        eend $?
}
