# admin/crontab/openexam-audit.sh.in
#!/bin/sh
# 
# Audit data cleanup script. This script is typical run periodical as a cron job.
#
# Author: Anders Lövgren
# Date:   2016-04-27

# This script directory:
cwd="$(dirname $(realpath $0))"

# The phalcon directory:
root="$(realpath $cwd/../../phalcon-mvc)"

# Lock file for cron job:
lock="$root/cache/audit.lock"

# Web server user:
user="apache"

# Script to run:
script="$root/script/openexam.php"

# Extra options:
#options="--verbose --force"

if ! [ -e $script ]; then
  echo "$0: missing script $script"
  exit 1
fi
if ! [ -d $root ]; then
  echo "$0: missing application directory $root"
  exit 1
fi
if [ -e $lock ]; then
  echo "$0: lock file exist (previous job not finished?)"
  exit 1
fi

# Create lock file:
touch $lock

# Cleanup expired audit data:
/usr/bin/sudo -u $user /usr/bin/php $script --audit --clean --days=90 $options $*

# Cleanup lock file:
rm -f $lock
